# bot.service
[Unit]
Description=ShortCircuit Bot Service
After=network.target

[Service]
# Replace with your actual username and path
User=ubuntu
WorkingDirectory=/home/ubuntu/ShortCircuit
ExecStart=/usr/bin/python3 /home/ubuntu/ShortCircuit/main.py
Restart=always
RestartSec=5

# Environment variables (Optional, better to use .env file in dir)
# Environment="TELEGRAM_BOT_TOKEN=xyz"

[Install]
WantedBy=multi-user.target
